# src/agent/kata-agent.service.in
[Unit]
Description=Kata Containers Agent
Documentation=https://github.com/kata-containers/kata-containers
Wants=kata-containers.target

[Service]
# Send agent output to tty to allow capture debug logs
# from a VM vsock port
StandardOutput=tty
Type=simple
# Load kernel modules before starting agent
# Run depmod first to ensure dependencies are up to date, then load erofs
# The - prefix means don't fail the service if the command fails
ExecStartPre=-/bin/bash -c '/sbin/depmod -a $(uname -r)'
ExecStartPre=-/sbin/modprobe erofs
ExecStart=@BINDIR@/@AGENT_NAME@
LimitNOFILE=1048576
# ExecStop is required for static agent tracing; in all other scenarios
# the runtime handles shutting down the VM.
ExecStop=/bin/sync ; /usr/bin/systemctl --force poweroff
FailureAction=poweroff
# Discourage OOM-killer from touching the agent
OOMScoreAdjust=-997